# systemd/feishu-bot@.service
[Unit]
Description=Feishu Bot (%i)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/feishu_bot/apps/%i
EnvironmentFile=-/root/feishu_bot/apps/%i/.env
ExecStart=/root/feishu_bot/apps/%i/.venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
